# configs/oram.service
[Unit]
Description=ORAM server
After=network.target
Documentation=https://pathoram.dbogatov.org/

[Service]
Type=simple

WorkingDirectory=/home/dmytro/dp-oram/dp-oram

ExecStartPre=/usr/bin/git pull
ExecStartPre=/usr/bin/make server

ExecStart=/home/dmytro/dp-oram/dp-oram/bin/oram-server

TimeoutStopSec=0
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
Alias=oram.service
